# systemd/wl_shimeji.socket
[Socket]
ListenUnix=/run/user/%i/shimeji-overlayd.sock
Accept=no

[Install]
WantedBy=sockets.target
